# dot_config/systemd/user/chezmoi-update.timer
[Unit]
Description=Update dotfiles daily

[Timer]
OnStartupSec=15m
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
